# kernel-logs.service
[Unit]
Description=Save kernel logs
DefaultDependencies=false
After=local-fs.target systemd-journald.service

[Service]
Type=oneshot
ExecStart=/etc/kernel-logs/kernel-logs.sh
StandardOutput=journal
StandardError=journal
RemainAfterExit=no

[Install]
WantedBy=sysinit.target
